# deployment/roles/pulseaudio/templates/pulseaudio.service.j2
[Unit]
Description=PulseAudio Sound Server for Open Cinema
Documentation=man:pulseaudio(1)
After=sound.target network.target

[Service]
Type=notify
User={{ open_cinema.user }}
Group={{ open_cinema.group }}
SupplementaryGroups=audio video

# Use user mode (not system mode)
ExecStart=/usr/bin/pulseaudio --daemonize=no --log-target=journal --exit-idle-time=-1
Restart=on-failure
RestartSec=5

# Runtime directory for PulseAudio socket
RuntimeDirectory=pulse
RuntimeDirectoryMode=0700

# Environment
Environment="PULSE_RUNTIME_PATH=/run/pulse"

# Security settings - more permissive for audio hardware access
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
