# deployment/roles/camilladsp/templates/camilladsp.service.j2
[Unit]
Description=CamillaDSP Audio Processor
Documentation=https://github.com/HEnquist/camilladsp
After=network.target sound.target pulseaudio.service

[Service]
Type=simple
User={{ open_cinema.user }}
Group={{ open_cinema.group }}
SupplementaryGroups=audio video

# CamillaDSP command with websocket enabled
ExecStart={{ camilladsp.install_path }} -p {{ camilladsp.websocket_port }} -a {{ camilladsp.websocket_address }} {{ camilladsp.config_path }}

# Restart policy
Restart=always
RestartSec=5

# Environment - PulseAudio socket
Environment="PULSE_RUNTIME_PATH=/run/pulse"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
